# ansible/roles/postgresql/templates/pgbackrest-full.timer.j2
[Unit]
Description=pgBackRest weekly full backup

[Timer]
OnCalendar=Sun *-*-* 19:00:00
Persistent=true

[Install]
WantedBy=timers.target
