# ansible/roles/postgresql/templates/pgbackrest-full.service.j2
[Unit]
Description=pgBackRest full backup
After=postgresql.service

[Service]
Type=oneshot
User=postgres
ExecStart=pgbackrest --stanza=main --type=full backup
