# ansible/roles/postgresql/templates/pgbackrest-backup.timer.j2
[Unit]
Description=pgBackRest daily differential backup

[Timer]
OnCalendar=*-*-* 19:30:00
Persistent=true

[Install]
WantedBy=timers.target
