# ansible/roles/postgresql/templates/pgbackrest-backup.service.j2
[Unit]
Description=pgBackRest differential backup
After=postgresql.service

[Service]
Type=oneshot
User=postgres
ExecStart=pgbackrest --stanza=main --type=diff backup
